# contrib/systemd/ioserver@.service
[Unit]
Description=IoSocks Server Service
After=network.target
Documentation=man:iosocks(8)

[Service]
Type=forking
ExecStart=/usr/bin/ioserver -d -c /etc/iosocks/%i.conf --pidfile /run/ioserver.pid --logfile /var/log/ioserver.log
PIDFile=/run/ioserver.pid

[Install]
WantedBy=multi-user.target